# supervisor/airSuck-federatedSysConn.conf
[program:airSuck-federatedSysConn]
directory=/opt/airSuck
command=python /opt/airSuck/federatedSysConn.py
autostart=true
autorestart=true
user=nobody
