# supervisor/airSuck-airSuckClient.conf
[program:airSuck-airSuckClient]
directory=/opt/airSuck
command=python /opt/airSuck/airSuckClient.py
autostart=true
autorestart=true
