# supervisor/airSuck-aisStateEngine.conf
[program:airSuck-aisStateEngine]
directory=/opt/airSuck
command=python /opt/airSuck/aisStateEngine.py
autostart=true
autorestart=true
