# supervisor/airSuck-dump1090ConnClt.conf
[program:airSuck-dump1090ConnClt]
directory=/opt/airSuck
command=python /opt/airSuck/dump1090ConnClt.py
autostart=true
autorestart=true
user=nobody