# supervisor/airSuck-stateNode.conf
[program:airSuck-stateNode]
directory=/opt/airSuck/node
command=nodejs /opt/airSuck/node/airSuck-stateNode.js
autostart=true
autorestart=true
