# supervisor/airSuck-ssrStateEngine.conf
[program:airSuck-ssrStateEngine]
directory=/opt/airSuck
command=python /opt/airSuck/ssrStateEngine.py
autostart=true
autorestart=true
user=nobody