# supervisor/airSuck-mongoDump.conf
[program:airSuck-mongoDump]
directory=/opt/airSuck
command=python /opt/airSuck/mongoDump.py
autostart=true
autorestart=true

