# supervisor/airSuck-dump1090Client.conf
[program:airSuck-dump1090Client]
directory=/opt/airSuck/node
command=nodejs /opt/airSuck/node/dump1090Client.js
autostart=true
autorestart=true
