# supervisor/airSuck-dump1090Connector.conf
[program:airSuck-dump1090Connector]
directory=/opt/airSuck
command=python /opt/airSuck/dump1090Connector.py
autostart=true
autorestart=true
