# supervisor/airSuck-dump1090ConnSrv.conf
[program:airSuck-dump1090ConnSrv]
directory=/opt/airSuck
command=python /opt/airSuck/dump1090ConnSrv.py
autostart=true
autorestart=true
