# supervisor/airSuck-airSuckServer.conf
[program:airSuck-airSuckServer]
directory=/opt/airSuck
command=python /opt/airSuck/airSuckServer.py
autostart=true
autorestart=true
