# supervisor/airSuck-stateMongoDump.conf
[program:airSuck-stateMongoDump]
directory=/opt/airSuck
command=python /opt/airSuck/stateMongoDump.py
autostart=true
autorestart=true
user=nobody